# ml-logic-service/supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=/usr/sbin/cron -f
autostart=true
autorestart=true
startretries=3

# [program:webserver]
# command=python3 /opt/src/server.py